# contrib/init/civilbitd.conf
description "Civilbit Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CIVILBITD_BIN="/usr/bin/civilbitd"
env CIVILBITD_USER="civilbit"
env CIVILBITD_GROUP="civilbit"
env CIVILBITD_PIDDIR="/var/run/civilbitd"
# upstart can't handle variables constructed with other variables
env CIVILBITD_PIDFILE="/var/run/civilbitd/civilbitd.pid"
env CIVILBITD_CONFIGFILE="/etc/civilbit/civilbit.conf"
env CIVILBITD_DATADIR="/var/lib/civilbitd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # civilbitd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CIVILBITD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run civilbitd."
        echo "The setting must appear in $CIVILBITD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Civilbit Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CIVILBITD_PIDDIR"
    chmod 0755 "$CIVILBITD_PIDDIR"
    chown $CIVILBITD_USER:$CIVILBITD_GROUP "$CIVILBITD_PIDDIR"
    chown $CIVILBITD_USER:$CIVILBITD_GROUP "$CIVILBITD_CONFIGFILE"
    chmod 0660 "$CIVILBITD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CIVILBITD_PIDFILE" \
    --chuid $CIVILBITD_USER:$CIVILBITD_GROUP \
    --exec "$CIVILBITD_BIN" \
    -- \
    -pid="$CIVILBITD_PIDFILE" \
    -conf="$CIVILBITD_CONFIGFILE" \
    -datadir="$CIVILBITD_DATADIR" \
    -disablewallet \
    -daemon

